# debian/ergo/DEBIAN/prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/ergo" ]; then
	invoke-rc.d ergo stop || exit $?
fi
# End automatically added section
